# scripts/fastintercom-sync.service
[Unit]
Description=FastIntercom Sync Daemon - Maintains fresh Intercom data
Documentation=https://github.com/yourusername/ask-intercom
After=network.target

[Service]
Type=simple
# Update with your actual paths
User=youruser
WorkingDirectory=/home/youruser/ask-intercom
Environment="PATH=/home/youruser/.local/bin:/usr/local/bin:/usr/bin:/bin"
Environment="INTERCOM_ACCESS_TOKEN=your_token_here"

# Run the sync daemon
ExecStart=/home/youruser/.local/bin/poetry run python scripts/fastintercom-sync-daemon.py --interval 15

# Restart on failure
Restart=on-failure
RestartSec=60

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=fastintercom-sync

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/youruser/.fastintercom

[Install]
WantedBy=multi-user.target
